# share/systemd/user/archiveserverd.service
[Unit]
Description=Video archive server

[Service]
ExecStart=/home/danila/.local/bin/video_archive_server

[Install]
WantedBy=default.target

