# share/systemd/user/streamrecorderd.service
[Unit]
Description=Web-camera recording daemon

[Service]
ExecStart=/home/danila/.local/bin/recorder

[Install]
WantedBy=default.target

