# share/systemd/user/videoremoverd.service
[Unit]
Description=Web-camera recording daemon

[Service]
ExecStart=/home/danila/.local/bin/remover

[Install]
WantedBy=default.target

